# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
autorestart=false
redirect_stderr=true

[program:php7-fpm]
command=/usr/sbin/php-fpm7.0
autorestart=false
redirect_stderr=true
